# bt-mqtt-gateway/rootfs/etc/services.d/btmqttgw/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Home Assistant Community Add-on: BT-MQTT Gateway
# Runs BT-MQTT Gateway
# ==============================================================================

bashio::log.info "Starting BT-MQTT Gateway..."

setcap cap_net_raw+eip $(eval readlink -f `which node`)

bashio::log.info "Installing default packages..."
pip3 install paho-mqtt pyyaml interruptingcow apscheduler tenacity

system_packages=$(bashio::config 'system_packages')
bashio::log.info "Installing $system_packages"
pip3 install $system_packages

bashio::log.info "Initilizing BT-MQTT Gateway..."
python3 /etc/btmqttgw/gateway.py
